# dist/systemd/birdwatcher.service
[Unit]
Description=healthchecker for BIRD-anycasted services
Documentation=https://github.com/skoef/birdwatcher
After=bird.service
After=bird6.service

[Service]
Type=simple
Environment=CONFIG_FILE=/etc/birdwatcher.conf
EnvironmentFile=-/etc/default/birdwatcher
ExecStartPre=/usr/bin/birdwatcher -config $CONFIG_FILE -check-config
ExecStart=/usr/bin/birdwatcher -config $CONFIG_FILE
Restart=on-failure

[Install]
WantedBy=multi-user.target
